# autostart/diodon.desktop
[Desktop Entry]
Type=Application
Exec=diodon
Hidden=false
X-MATE-Autostart-enabled=true
Name[en_US]=diodon
Name=diodon
Comment[en_US]=
Comment=
